# files/supervisord.conf
[supervisord]
nodaemon=true

[program:syncthing]
user = root
autostart = true
stdout_logfile = /var/log/syncthing-stdout.log
stderr_logfile = /var/log/syncthing-stderr.log
